# installer/config/backend.service
[Unit]
Description=Backend FastAPI Service
After=network.target

[Service]
User=beer
Group=www-data
WorkingDirectory= /var/www/Testttt/apiii
Environment="PATH=/usr/local/Testttt/bin"
ExecStart=/usr/bin/python3 -m uvicorn apiii.main:app --host 192.168.21.132 --port 8080
Restart=always

[Install]
WantedBy=multi-user.target
